# systemd/market-scraper.service
[Unit]
Description=CryptoData Market Scraper Service
Documentation=https://github.com/monjurkuet/cryptoOS
After=network.target mongodb.service redis.service
Wants=mongodb.service redis.service

[Service]
Type=simple
User=administrator
Group=administrator
WorkingDirectory=/home/administrator/githubrepo/cryptoOS/market-scraper
Environment="PATH=/home/administrator/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="UV_CACHE_DIR=/home/administrator/githubrepo/cryptoOS/logs/.uv-cache"
Environment="LOG_LEVEL=INFO"
ExecStart=/home/administrator/.local/bin/uv run python -m market_scraper server --host 127.0.0.1 --port 3845
Restart=always
RestartSec=5
StandardOutput=append:/home/administrator/githubrepo/cryptoOS/logs/market-scraper.log
StandardError=append:/home/administrator/githubrepo/cryptoOS/logs/market-scraper.error.log

# Resource limits - optimized for budget constraints
MemoryMax=1G
CPUQuota=80%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/administrator/githubrepo/cryptoOS/logs

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
