# systemd/signal-system.service
[Unit]
Description=CryptoData Smart Money Signal System
Documentation=https://github.com/monjurkuet/cryptoOS
After=network.target redis.service market-scraper.service
Wants=redis.service
Requires=market-scraper.service

[Service]
Type=simple
User=administrator
Group=administrator
WorkingDirectory=/home/administrator/githubrepo/cryptoOS/smart-money-signal-system
Environment="PATH=/home/administrator/githubrepo/cryptoOS/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="LOG_LEVEL=INFO"
Environment="API_HOST=127.0.0.1"
Environment="API_PORT=4341"
ExecStart=/home/administrator/githubrepo/cryptoOS/.venv/bin/python -m signal_system server
Restart=always
RestartSec=5
StandardOutput=append:/home/administrator/githubrepo/cryptoOS/logs/signal-system.log
StandardError=append:/home/administrator/githubrepo/cryptoOS/logs/signal-system.error.log

# Resource limits
MemoryMax=512M
CPUQuota=40%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/administrator/githubrepo/cryptoOS/logs

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
